# bin/install/amavis-stats-0.1.22/startup/amavisstats.in
#!/sbin/runscript

depend() {
        need net
        use logger
        after amavisd
}

start() {
        ebegin "Starting amavis-stats"
        start-stop-daemon --start --quiet --exec @SBINDIR@/@PKGNAME@
        eend $?
}

stop() {
        ebegin "Stopping amavisd-stats"
        start-stop-daemon --stop --quiet --pidfile @PKGLIBDIR@/@PKGNAME@/@PKGNAME@.pid
        eend $?
}

